# netmonitor.service
[Unit]
Description=Network Monitor - Verdacht Verkeer Detectie
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netmonitor
ExecStart=/usr/bin/python3 /opt/netmonitor/netmonitor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
