# netmonitor-feed-update.service
# NOTE: This is a TEMPLATE file. Do not use directly!
# Use install_services.sh which will generate the correct service file
# with the actual paths for your installation.
#
# Example paths shown below - will be replaced by install script:

[Unit]
Description=Network Monitor - Threat Feed Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/path/to/netmonitor
ExecStart=/path/to/netmonitor/venv/bin/python3 /path/to/netmonitor/update_feeds.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
