# services/netmonitor-mcp-sse.service.template
[Unit]
Description=NetMonitor MCP SSE Server (Legacy - Claude Desktop)
Documentation=https://github.com/willempoort/netmonitor
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__/mcp_server/legacy_stdio_sse

# Set defaults first - these can be overridden by EnvironmentFile
Environment="MCP_SSE_HOST=0.0.0.0"
Environment="MCP_SSE_PORT=3000"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_DIR=/var/log/netmonitor"

# Load environment variables from .env (overrides defaults above)
EnvironmentFile=-__INSTALL_DIR__/.env
EnvironmentFile=-/etc/netmonitor/netmonitor.env

# Start MCP SSE server (legacy transport for Claude Desktop)
ExecStart=__INSTALL_DIR__/venv/bin/python3 __INSTALL_DIR__/mcp_server/legacy_stdio_sse/server.py \
    --transport sse \
    --host $MCP_SSE_HOST \
    --port $MCP_SSE_PORT

# Restart behavior
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/netmonitor /var/lib/netmonitor __INSTALL_DIR__

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netmonitor-mcp-sse

[Install]
WantedBy=multi-user.target
