# services/netmonitor-mcp-streamable.service.template
[Unit]
Description=NetMonitor MCP Streamable HTTP Server (Claude Desktop & Open-WebUI)
Documentation=https://github.com/willempoort/netmonitor
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__

# Set defaults first - these can be overridden by EnvironmentFile
Environment="MCP_API_HOST=127.0.0.1"
Environment="MCP_API_PORT=8000"
Environment="MCP_DEBUG=false"
Environment="LOG_LEVEL=INFO"
Environment="PYTHONUNBUFFERED=1"

# Load environment variables from .env (overrides defaults above)
EnvironmentFile=-__INSTALL_DIR__/.env
EnvironmentFile=-/etc/netmonitor/netmonitor.env

# Start MCP Streamable HTTP server
ExecStart=__INSTALL_DIR__/venv/bin/python3 __INSTALL_DIR__/mcp_server/streamable_http_server.py

# Restart behavior
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/netmonitor /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netmonitor-mcp-streamable

[Install]
WantedBy=multi-user.target
